# launchd/com.ai-bridge.mlx.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ai-bridge.mlx</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>source ~/.ai-bridge/mlx-venv/bin/activate && python -m mlx_lm.server --host 127.0.0.1 --port 8080 --model mlx-community/Qwen2.5-Coder-7B-Instruct-4bit</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>~/.ai-bridge/logs/mlx.out</string>

    <key>StandardErrorPath</key>
    <string>~/.ai-bridge/logs/mlx.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
